# hardware-lab-v4/infra/systemd/reboot.timer
[Unit]
Description=Reboot diario as 04:00

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
