# hardware-lab-v4/infra/systemd/reboot.service
[Unit]
Description=Reboot diario do sistema

[Service]
Type=oneshot
ExecStart=/sbin/reboot
